# glue-validator-cron/tags/R_1_3_0_1/etc/cron.hourly/glue-validator-cron
#!/bin/sh

glue-validator -H localhost -p 2170 -b GLUE2GroupID=resource,o=glue -n -v 3 &> /var/log/glue-validator/glue2-validator.log

